# scripts/RR_to_SS_uriage-arari/rr-uriage-arari.service
[Unit]
Description=RR to SS 売上・粗利実績同期
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/project/apiconnection_factory/scripts/RR_to_SS_uriage-arari
EnvironmentFile=/etc/apiconnection.env
ExecStart=/root/project/apiconnection_factory/venv/bin/python RR_to_SS_uriage-arari.py
StandardOutput=journal
StandardError=journal
