# scripts/RR_to_RMS_Lstage-customerdata/rms-to-rr-shippingday-reflect.timer
[Unit]
Description=RMS出荷日→楽楽販売反映処理 タイマー (毎日9時)

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
