# scripts/RR_to_RMS_Lstage-customerdata/rr-customerdata.timer
[Unit]
Description=RR to RMS Lstage 顧客データ同期 タイマー (30分毎)

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
