# scripts/RR_to_SS_seiyakuprocess/rr-seiyakuprocess.timer
[Unit]
Description=RR to SS 成約プロセス同期 タイマー (1時間ごと)

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
