# scripts/RR_to_SS_UUcounter/rr-uucounter.timer
[Unit]
Description=RR to SS UUカウンター同期 タイマー (毎日8時)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
